# debian/kwtools-net-narc.postrm
#! /bin/sh
# postrm script for kwtools-net-narc
#
# see: dh_installdeb(1)

set -e

if [ -x "`which update-menus 2>/dev/null`" ] ; then
	update-menus
fi

case "$1" in
	purge)
		if [ -e /usr/share/debconf/confmodule ] ; then
			. /usr/share/debconf/confmodule
			db_purge
		fi
		update-rc.d -f kwnarc remove >/dev/null || exit $?
		;;
esac

#DEBHELPER#

exit 0
